# src/drivers/usb/x10_type2/driver.rules
ACTION!="add", GOTO="liri_end"

ENV{ID_VENDOR_ID}=="0bc7", ENV{ID_MODEL_ID}=="0006", ENV{liri_receiver_driver}="x10_type2", ENV{liri_receiver_layout}="x10.0536", RUN+="/usr/bin/liricallout"

LABEL="liri_end"
